# message-processor-cron.service
[Unit]
Description=Message Processor Cron Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=your_username
WorkingDirectory=~/mnt/narayan/Code/ambient_agents/git_slack_agent
Environment=PATH=~/mnt/narayan/Code/ambient_agents/git_slack_agent/langgraph_env/bin
ExecStart=~/mnt/narayan/Code/ambient_agents/git_slack_agent/langgraph_env/bin/python ~/mnt/narayan/Code/ambient_agents/git_slack_agent/cron_message_processor.py --interval 30
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Logging
LogLevel=info

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=~/mnt/narayan/Code/ambient_agents/git_slack_agent

[Install]
WantedBy=multi-user.target 